# dht_mqtt.service
[Unit]
Description=DHT MQTT daemon

[Service]
Type=simple
ExecStart=%{python} %{filepath}

[Install]
WantedBy=multi-user.target